# packaging/bilbycast-edge.service
[Unit]
Description=bilbycast-edge media transport gateway
Documentation=https://bilbycast.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bilbycast
Group=bilbycast

# `current` is a symlink the upgrade pipeline atomically swaps; systemd
# follows it on every (re)start so a successful staged upgrade lands the
# moment the old binary exits.
ExecStart=/opt/bilbycast/edge/current/bilbycast-edge --config /opt/bilbycast/edge/config.json

# Restart on any non-zero exit so the boot watchdog can drive the
# rollback handshake. 3-second back-off is plenty — typical respawn
# latency is dominated by systemd's own settling, not the binary.
Restart=always
RestartSec=3
# Limit how often we'll respawn; pairs with the watchdog's
# `max_boot_attempts` so a crash-loop on a bad upgrade doesn't burn
# CPU forever before the symlink revert kicks in.
StartLimitInterval=300
StartLimitBurst=10

# Send SIGTERM and let the edge drain flows + close sockets cleanly.
KillSignal=SIGTERM
TimeoutStopSec=20

# ── Hardening ──
ProtectSystem=strict
# Writable paths: the install root (versions/, current symlink,
# state.json), the data dir (replay segments, media library, etc.),
# and the runtime tmp dir. Everything else is read-only.
ReadWritePaths=/opt/bilbycast/edge /var/lib/bilbycast/edge /run/bilbycast
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
# Allow opening DRI render nodes (VAAPI, KMS scanout for `display`)
# and ALSA devices when those features are compiled in. Operators
# who don't need the display output can drop these without affecting
# the rest of the edge.
DeviceAllow=/dev/dri/card0 rwm
DeviceAllow=/dev/dri/card1 rwm
DeviceAllow=/dev/dri/renderD128 rwm
DeviceAllow=/dev/dri/renderD129 rwm
DeviceAllow=/dev/snd/* rwm
SupplementaryGroups=video render audio
# Capability set. We need exactly one capability:
#
#   CAP_NET_ADMIN — required by mainline kernel ≥ 6.x (and every recent
#   Ubuntu / RHEL backport) for `setsockopt(SO_TXTIME)` with any non-
#   `CLOCK_MONOTONIC` clockid. The wire-pacing path uses `CLOCK_TAI`
#   (the only clockid the kernel ETF qdisc on Intel ice/igc/igb and
#   Mellanox mlx5 will accept, and the natural choice when ptp4l +
#   phc2sys are running). Without this cap, `enable_so_txtime` returns
#   EPERM and `engine::wire_emit` silently falls back to the
#   `clock_nanosleep` tier — degraded pacing precision (≈ 50–500 µs
#   jitter vs sub-µs with NIC offload), and on any host where the ETF
#   qdisc has `skip_sock_check off`, EVERY non-SO_TXTIME packet is
#   also dropped at the qdisc — so the fallback path silently fails
#   completely.
#
# CAP_NET_ADMIN does NOT let the edge install qdiscs (we still want
# that to be operator-side at provision time via `setup-etf-qdisc.sh`)
# — it only unlocks the per-socket setsockopt path.
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
RestrictNamespaces=true
# `wire_emit` uses SCHED_FIFO best-effort on its egress thread for
# sub-100 µs `clock_nanosleep` precision (the alternative is the
# tokio timer wheel's ~1–5 ms ceiling, which is too coarse for
# PCR-anchored pacing). RestrictRealtime=true would block this. Pair
# with `LimitRTPRIO=50` below — the kernel allows unprivileged
# `sched_setscheduler(SCHED_FIFO)` whenever RTPRIO is non-zero, so
# we don't need to grant `CAP_SYS_NICE`. Wire_emit logs at debug and
# falls back to default scheduler if the priority change is denied,
# so there is no hard failure path.
RestrictRealtime=false
LockPersonality=true
MemoryDenyWriteExecute=false  # JIT crates exist (rustls, ring) — leave on default
SystemCallArchitectures=native

# Resource limits — generous defaults; operators can tune via drop-in
# /etc/systemd/system/bilbycast-edge.service.d/override.conf
LimitNOFILE=65536
# `LimitMEMLOCK=infinity` together with `BILBYCAST_MLOCKALL=1` in the
# environment file below makes the edge pin every current + future
# page into RAM at startup via `mlockall(MCL_CURRENT | MCL_FUTURE)`.
# Eliminates major-page-fault stalls on the data-plane hot path —
# without this, a single page fault on a wire-emit thread can stall
# pacing for 10–100 ms and drop packets at the receiver.
LimitMEMLOCK=infinity
# RT-priority ceiling. Wire-emit requests SCHED_FIFO priority 50.
# Headroom up to 99 lets future hot-path threads (codec pool at 40,
# PID-bus assembler at 30, PCR PLL sampler at 35 — see
# docs/production-tuning.md) acquire SCHED_FIFO without
# `CAP_SYS_NICE`. Default RLIMIT_RTPRIO is 0 on most distros, which
# would force a silent fallback to SCHED_OTHER (1–5 ms p99 pacing
# jitter — too coarse for broadcast PCR_AC).
LimitRTPRIO=99

# Environment file lets operators set RUST_LOG, BILBYCAST_REPLAY_DIR,
# BILBYCAST_MLOCKALL=1, BILBYCAST_WIRE_EMIT_CPUS=<set>, etc. without
# editing the unit. Created (with `BILBYCAST_MLOCKALL=1` seeded) by
# install-edge.sh.
EnvironmentFile=-/etc/bilbycast/edge.env

[Install]
WantedBy=multi-user.target
