# packaging/bilbycast-edge.service
[Unit]
Description=bilbycast-edge media transport gateway
Documentation=https://bilbycast.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bilbycast
Group=bilbycast

# `current` is a symlink the upgrade pipeline atomically swaps; systemd
# follows it on every (re)start so a successful staged upgrade lands the
# moment the old binary exits.
ExecStart=/opt/bilbycast/edge/current/bilbycast-edge --config /opt/bilbycast/edge/config.json

# Restart on any non-zero exit so the boot watchdog can drive the
# rollback handshake. 3-second back-off is plenty — typical respawn
# latency is dominated by systemd's own settling, not the binary.
Restart=always
RestartSec=3
# Limit how often we'll respawn; pairs with the watchdog's
# `max_boot_attempts` so a crash-loop on a bad upgrade doesn't burn
# CPU forever before the symlink revert kicks in.
StartLimitInterval=300
StartLimitBurst=10

# Send SIGTERM and let the edge drain flows + close sockets cleanly.
KillSignal=SIGTERM
TimeoutStopSec=20

# ── Hardening ──
NoNewPrivileges=true
# Allow opening DRI render nodes (VAAPI, KMS scanout for `display`)
# and ALSA devices when those features are compiled in. Class-based
# (`char-<name>` per /proc/devices) rather than node paths: covers any
# number of GPUs / sound cards and survives hotplug renumbering, where
# a hardcoded /dev/dri/cardN list breaks on the third GPU.
DeviceAllow=char-drm rwm
DeviceAllow=char-alsa rwm
SupplementaryGroups=video render audio
# CAP_NET_ADMIN — required by kernel ≥ 6.x for setsockopt(SO_TXTIME)
# with CLOCK_TAI. Without it, wire-emit falls back to clock_nanosleep.
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
RestrictRealtime=false
SystemCallArchitectures=native

# Resource limits — generous defaults; operators can tune via drop-in
# /etc/systemd/system/bilbycast-edge.service.d/override.conf
LimitNOFILE=65536
# `LimitMEMLOCK=infinity` together with `BILBYCAST_MLOCKALL=1` in the
# environment file below makes the edge pin every current + future
# page into RAM at startup via `mlockall(MCL_CURRENT | MCL_FUTURE)`.
# Eliminates major-page-fault stalls on the data-plane hot path —
# without this, a single page fault on a wire-emit thread can stall
# pacing for 10–100 ms and drop packets at the receiver.
LimitMEMLOCK=infinity
# RT-priority ceiling. Wire-emit requests SCHED_FIFO priority 50.
# Headroom up to 99 lets future hot-path threads (codec pool at 40,
# PID-bus assembler at 30, PCR PLL sampler at 35 — see
# docs/production-tuning.md) acquire SCHED_FIFO without
# `CAP_SYS_NICE`. Default RLIMIT_RTPRIO is 0 on most distros, which
# would force a silent fallback to SCHED_OTHER (1–5 ms p99 pacing
# jitter — too coarse for broadcast PCR_AC).
LimitRTPRIO=99

# Environment file lets operators set RUST_LOG, BILBYCAST_REPLAY_DIR,
# BILBYCAST_MLOCKALL=1, BILBYCAST_WIRE_EMIT_CPUS=<set>, etc. without
# editing the unit. Created (with `BILBYCAST_MLOCKALL=1` seeded) by
# install-edge.sh.
EnvironmentFile=-/etc/bilbycast/edge.env

[Install]
WantedBy=multi-user.target
